# .config/systemd/user/whisper-dictation.service
[Unit]
Description=Local Whisper server for system-wide dictation

[Service]
Type=simple
ExecStart=%h/.local/opt/whisper.cpp/build/bin/whisper-server --host 127.0.0.1 --port 8178 --model %h/.local/share/whisper/ggml-base.en.bin --threads 4 --language en --no-gpu --no-timestamps
WorkingDirectory=%h/.local/opt/whisper.cpp
Restart=on-failure
RestartSec=2

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_UNIX AF_INET
